# sonarr/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Sonarr
# Runs the Nginx daemon
# ==============================================================================

# Wait for Sonarr to become available
bashio::net.wait_for 8989 localhost 300

bashio::log.info "Starting NGinX..."
exec nginx
