# sabnzbd/rootfs/etc/services.d/sabnzbd/run
#!/usr/bin/with-contenv bash

exec \
  s6-setuidgid abc python3 /app/sabnzbd/SABnzbd.py \
  --config-file /data/config --server 0.0.0.0:8080 --browser 0
